# deployment/systemd/nexus-scheduler.service
[Unit]
Description=Nexus Monitoring Portal Celery Beat Scheduler
Documentation=https://github.com/your-repo/nexus-monitoring-portal
After=network.target redis.service postgresql.service nexus-worker.service
Requires=redis.service nexus-worker.service
Wants=postgresql.service

[Service]
Type=simple
User=nexus
Group=nexus
WorkingDirectory=/home/nexus/monitoring-portal/backend
Environment=PATH=/home/nexus/monitoring-portal/backend/venv/bin
EnvironmentFile=/home/nexus/monitoring-portal/backend/.env.production

# Main process
ExecStart=/home/nexus/monitoring-portal/backend/venv/bin/celery -A tasks.celery_app beat --loglevel=info --schedule=/var/lib/nexus/celerybeat-schedule

# Process management
Restart=always
RestartSec=3
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/nexus/monitoring-portal
ReadWritePaths=/var/log/nexus
ReadWritePaths=/var/lib/nexus

# Resource limits
LimitNOFILE=65536
LimitNPROC=1024

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nexus-scheduler

# Health check
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target